# HW8_bash/parser.service
[Unit]
Description="Log analysis"

[Service]
Type=simple
ExecStart=/vagrant/parser.sh 4 3

