# HW7_systemd/monitoring/monitoring.service
[Unit]
Description="Monitoring for keyword"

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/monitoring
ExecStart=/opt/monitoring.sh 

