# 08_bash/parser.service
[Unit]
Description="Log analysis"

[Service]
Type=simple
ExecStart=/bin/bash /vagrant/parser.sh 4 3

