# project/roles/wordpress/templates/borg-backup.timer.j2
[Unit]
Description=Borg backup timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=multi-user.target

